# com.bluetti.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bluetti.bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>VENV_PATH/bin/python</string>
        <string>SCRIPT_PATH/bluetti_bridge.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>SCRIPT_PATH</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>SCRIPT_PATH/logs/bridge.log</string>

    <key>StandardErrorPath</key>
    <string>SCRIPT_PATH/logs/bridge.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>VENV_PATH/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
